# inventory_app/service/inventory-backup.timer
[Unit]
Description=Daily Inventory DB backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
